# Source/Host/openWRT/backfire/package/6LoWPAN/files/etc/uci-defaults/luci-theme-nxp
#!/bin/sh
	uci batch <<-EOF
	set luci.themes.NXP=/luci-static/nxp
	set luci.main.mediaurlbase=/luci-static/nxp
	commit luci
EOF
