# Source/Host/openWRT/targets/lpc32xx/backfire-changes/target/linux/lpc32xx/base-files/etc/uci-defaults/dr1182
#!/bin/sh
#
# Copyright (C) 2010 OpenWrt.org
#

. /lib/lpc32xx.sh

board=$(lpc32xx_board_name)

dr1182_set_leds() {
	uci batch <<EOF
add system led
set system.@led[-1]='led'
set system.@led[-1].name='Heartbeat'
set system.@led[-1].sysfs='oem:red:led1'
set system.@led[-1].trigger='heartbeat'
add system led
set system.@led[-1]='led'
set system.@led[-1].name='6LoWPAN Activity'
set system.@led[-1].sysfs='dr1182:green:led2'
set system.@led[-1].trigger='netdev'
set system.@led[-1].dev='tun0'
set system.@led[-1].mode='link tx rx'
add system led
set system.@led[-1]='led'
set system.@led[-1].name='Ethernet Activity'
set system.@led[-1].sysfs='dr1182:green:led3'
set system.@led[-1].trigger='netdev'
set system.@led[-1].dev='eth0'
set system.@led[-1].mode='link tx rx'
commit system
EOF
}

if [ "${board}" == "dr1182" ]; then
	dr1182_set_leds
fi
